# crontab
SHELL=/bin/sh

0 * * * * php /app/bin/app.php sessions:cleanup